# deploy/mx-server/postfix_label_proxy.service
# ──────────────────────────────────────────────────────────────────────────────
# UNIT 2/2: prom-label-proxy
# Injeta a label "cliente" em todas as métricas e expõe na porta 9155.
# O OTel Collector centralizado deve raspar ESTA porta.
#
# Instalação do binário:
#   VERSAO=0.11.0
#   curl -L "https://github.com/prometheus-community/prom-label-proxy/releases/download/v${VERSAO}/prom-label-proxy_${VERSAO}_linux_amd64.tar.gz" \
#     | tar xz && mv prom-label-proxy /usr/local/bin/
# ──────────────────────────────────────────────────────────────────────────────
[Unit]
Description=Prometheus Label Proxy para postfix_exporter
After=network.target postfix_exporter.service
Requires=postfix_exporter.service

[Service]
Type=simple
User=nobody
Group=nogroup
# Substitua "nome-do-cliente" pelo valor real antes de instalar
Environment="CLIENTE=nome-do-cliente"
ExecStart=/usr/local/bin/prom-label-proxy \
    --insecure-listen-address=0.0.0.0:9155 \
    --upstream=http://127.0.0.1:9154 \
    --label=cliente=${CLIENTE}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
